# pod-mirror/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'

[program:flask]
command=python3 /home/generic/app.py
user=generic
