# pod-mirror/supervisord.conf
[supervisord]
user=generic
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
user=generic

[program:flask]
command=python3 /home/generic/app.py
user=generic
