# pod-mirror/supervisord.conf
[supervisord]
nodaemon=true
user=root 

[program:nginx]
command=nginx -g 'daemon off;'
user=root

[program:flask]
command=python3 /home/generic/app.py
user=generic
